# deploy/com.netmap.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.netmap.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/phil/Projects/NetMap/NetMapServer/.build/debug/App</string>
    </array>

    <!-- Restart automatically if the process exits -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Redirect stdout/stderr to the same log used by the VS Code tasks -->
    <key>StandardOutPath</key>
    <string>/tmp/netmap_server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/netmap_server.log</string>

    <!-- Working directory so relative paths (e.g. Public/) resolve correctly -->
    <key>WorkingDirectory</key>
    <string>/Users/phil/Projects/NetMap/NetMapServer</string>

    <!-- Throttle restarts: wait 5 s before restarting after a crash -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
